# senia-networking.install/etc/network/if-up.d/wake-up
#!/bin/bash

# Simple script to enable wakeonlan on if-up

logger -t senia-networking "Setting $IFACE enable to be wake-up"

ethtool -s $IFACE wol g || logger -t senia-networking "Setting $IFACE enable to be wake-up ... FAILED"

exit 0

